# .readme/root-config/crontab
*/3 * * * * ~/.config/i3/scripts/mail_fetch.sh
